# home/private_Library/LaunchAgents/io.aimuzov.session-env.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.aimuzov.session-env</string>
    <!--
      Прокидывает не-секретные session-переменные в launchd-сессию, чтобы их
      видели GUI-приложения и не-login shell'ы (например Bash-tool у Claude),
      которые не читают config.fish.

      GNUPGHOME: keyring лежит в XDG-пути, а не в дефолтном ~/.gnupg. Без этой
      переменной gpg в не-login окружении не находит секретный ключ и подпись
      коммитов падает с "No secret key". Дублирует config.fish осознанно — там
      она нужна для login-shell и Linux/VPS, где launchd нет.

      Секреты (GITHUB_TOKEN) сюда НЕ кладём — для них отдельный private-агент.
      Расширять списком: launchctl setenv FOO "bar"; launchctl setenv BAZ "qux"
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>launchctl setenv GNUPGHOME "{{ .chezmoi.homeDir }}/.local/share/gnupg"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
