# home/private_Library/LaunchAgents/private_io.aimuzov.session-env.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.aimuzov.session-env</string>
    <!--
      Прокидывает session-переменные в launchd-сессию, чтобы их видели
      GUI-приложения, запущенные из Dock/Finder (не через шелл) — они не читают
      ни config.fish, ни .zshenv.

      Файл private_ (права 600), т.к. среди переменных есть секреты из KeePassXC.

      Источник переменных — home/.chezmoitemplates/session-env.tmpl (режим launchctl).
      Добавить переменную = одна строка ТАМ; сюда она попадёт автоматически.
      Каждая команда на своей строке исполняется /bin/sh -c как отдельная
      (перевод строки = разделитель команд).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>{{ includeTemplate "session-env.tmpl" (dict "syntax" "launchctl" "homeDir" .chezmoi.homeDir) }}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
