# home/private_Library/LaunchAgents/private_io.aimuzov.session-env.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.aimuzov.session-env</string>
    <!--
      Pushes session variables into the launchd session so they are visible to
      GUI apps launched from Dock/Finder (not via a shell) — those read neither
      config.fish nor .zshenv.

      A private_ file (mode 600), since the variables include secrets from KeePassXC.

      Variable source — home/.chezmoitemplates/session-env.tmpl (launchctl mode).
      Adding a variable = one line THERE; it propagates here automatically.
      Each command on its own line is run by /bin/sh -c as a separate one
      (newline = command separator).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>{{ includeTemplate "session-env.tmpl" (dict "syntax" "launchctl" "homeDir" .chezmoi.homeDir) }}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
