# home/private_Library/LaunchAgents/private_io.aimuzov.session-env.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.aimuzov.session-env</string>
    <!--
      Прокидывает session-переменные в launchd-сессию, чтобы их видели
      GUI-приложения и не-login shell'ы (например Bash-tool у Claude), которые
      не читают config.fish.

      Файл private_ (права 600), т.к. содержит секрет GITHUB_TOKEN из KeePassXC.

        GNUPGHOME    — keyring в XDG-пути, а не дефолтном ~/.gnupg; без неё gpg
                       в не-login окружении не находит ключ ("No secret key").
                       Дублирует config.fish осознанно (нужна на Linux/VPS без launchd).
        GITHUB_TOKEN — токен для gh/github-MCP в GUI-приложениях.

      Новые НЕсекретные переменные — просто добавить "; launchctl setenv FOO bar".
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>launchctl setenv GNUPGHOME "{{ .chezmoi.homeDir }}/.local/share/gnupg"; launchctl setenv GITHUB_TOKEN "{{ (keepassxcAttribute "GITHUB" "TOKEN") }}"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
